# deploy/xelity-autoupdate.timer
[Unit]
Description=Xelity auto-update from GitHub (every 2 min)
After=network-online.target

[Timer]
# После завершения предыдущего update — ещё через 2 мин (oneshot)
OnBootSec=1min
OnUnitInactiveSec=2min
AccuracySec=15s
Persistent=true
Unit=xelity-autoupdate.service

[Install]
WantedBy=timers.target
